# resources/dashboard.sv.conf
[program:dashboard]
command=/bin/sh -c "npm --prefix /dashboard install /dashboard && node /dashboard/bin/www"